# packages/librenet6/files/common/etc/uci-defaults/tinc-setup
#!/bin/sh

HOSTNAME=$(uci get system.@system[0].hostname)

# tell tinc to use the librenet6 network
uci -q show tinc.librenet6 || {
    touch /etc/config/tinc
    uci batch <<-EOF
    set tinc.librenet6='tinc-net'
    set tinc.librenet6.enabled='1'
    commit tinc
EOF
}

# create keys if the don't exists yet
[[ -e /etc/tinc/librenet6/hosts/node_$HOSTNAME ]] || {
    if [[ -e /usr/sbin/tinc ]]; then
	# tinc 1.1 installed
	tinc -n librenet6 generate-keys </dev/null
	cat /etc/tinc/librenet6/*.pub > \
		/etc/tinc/librenet6/hosts/node_$HOSTNAME
    else
	# use old tincd
	tincd -n librenet6 -K </dev/null
   fi
}

# disable ra
uci set network.lan.accept_ra=0
uci commit network

